# systemd/geonext-mcp.service
[Unit]
After=network.target

[Service]
WorkingDirectory=/opt/geonext-mcp
EnvironmentFile=/opt/geonext-mcp/.env
ExecStart=/opt/geonext-mcp/.venv/bin/geonext-mcp
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
